# systemd-propagate/network-online.target
[Unit]
Description=Network is Online (propagated)
