# brightd/brightd.service
[Unit]
Description=Automatically fade screen when idle
After=systemd-udevd.service

[Service]
Type=simple
ExecStart=/usr/bin/brightd -vx
Restart=on-abort

[Install]
WantedBy=multi-user.target
